# scripts/picloud-node/supervisor/picloud-node-mqtt.conf
[program:picloud-node-mqtt]
command=/opt/picloud-node/stats.py
autostart=true
autorestart=true